# supervisord.conf
[supervisord]
nodaemon=true

[program:unbound]
command=/usr/sbin/unbound -c /etc/unbound/unbound.conf -d
